# .ddev/web-build/drupal.cron
*/10 * * * * IS_DDEV_PROJECT=true DDEV_PHP_VERSION=8.3 /var/www/html/vendor/bin/drush cron | tee -a /var/www/html/cron.log
