# config/server/roles/laggyluke.nodejs-app/templates/upstart-master.conf.j2
{% if nodejs_app.autostart %}
start on runlevel [2345]
stop on runlevel [016]
{% endif %}
